# keylist/keylist.timer
[Unit]
Description=Boot Keylist

[Timer]
OnBootSec=60
Unit=keylist.service

[Install]
WantedBy=multi-user.target
